# examples/bub.timer
[Unit]
Description=Nightly bub timer

[Timer]
OnCalendar=daily
AccuracySec=12h

[Install]
WantedBy=timers.target
